# systemd/nos-pfe.service
[Unit]
Description=NOS Packet Forwarding Engine
Documentation=https://github.com/nos-dev/nos
After=network.target

[Service]
Type=notify
ExecStart=/usr/local/sbin/nos-pfe --xdp-obj /usr/lib/nos/xdp_prog.o
Restart=on-failure
RestartSec=5s

RuntimeDirectory=nos
RuntimeDirectoryMode=0750

CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN

LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
