# systemd/nos-cli.service
[Unit]
Description=NOS CLI Shell
Documentation=https://github.com/nos-dev/nos
After=network.target nos-pfe.service
Wants=nos-pfe.service

[Service]
Type=simple
ExecStartPre=/bin/bash -c 'chown root:nos /run/nos 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'chmod 775 /run/nos 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'test -S /run/nos/pfe.sock && chown root:nos /run/nos/pfe.sock && chmod 660 /run/nos/pfe.sock || true'
ExecStart=/usr/local/bin/nos-cli
User=nos
Group=nos
Restart=on-failure
RestartSec=3s
AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN

StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes

[Install]
WantedBy=multi-user.target
