# systemd/nos-pfe.service
[Unit]
Description=NOS Packet Forwarding Engine
Documentation=https://github.com/nos-dev/nos
After=network.target

[Service]
Type=notify
ExecStart=/usr/local/sbin/nos-pfe --xdp-obj /usr/lib/nos/xdp_prog.o
ExecStartPost=/bin/bash -c 'sleep 0.5 && chown root:nos /run/nos && chmod 775 /run/nos && chown root:nos /run/nos/pfe.sock && chmod 660 /run/nos/pfe.sock'
Restart=on-failure
RestartSec=5s

RuntimeDirectory=nos
RuntimeDirectoryMode=0775
RuntimeDirectoryGroup=nos

CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN CAP_CHOWN
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN CAP_CHOWN

LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
